# openrc/bern-auto-shutdown
#!/sbin/openrc-run
command=/usr/local/sbin/bern-auto-shutdown
#command_args=""
command_background=true
pidfile=/var/run/bern-auto-shutdown.pid
name="bern auto shutdown"

